# services/remountd.socket.in
[Unit]
Description=remountd control socket

[Socket]
ListenStream=@REMOUNTD_SOCKET_PATH@
SocketUser=root
SocketGroup=remountd
SocketMode=0660
Accept=no

[Install]
WantedBy=sockets.target
